# macos/com.dotfiles.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dotfiles.sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-l</string>
    <string>-c</string>
    <string>$HOME/.dotfiles/sync/auto-sync.sh</string>
  </array>

  <!-- Run every 30 minutes -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <!-- Run on load -->
  <key>RunAtLoad</key>
  <false/>

  <!-- Log output -->
  <key>StandardOutPath</key>
  <string>/tmp/dotfiles-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/dotfiles-sync-error.log</string>

  <!-- Only run when on AC power (don't drain laptop battery) -->
  <key>ProcessType</key>
  <string>Background</string>

  <key>LowPriorityIO</key>
  <true/>

  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
