# macos/com.captures.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.captures.watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-l</string>
    <string>-c</string>
    <string>bash $HOME/.dotfiles/scripts/move-captures.sh</string>
  </array>

  <!-- Watch the Desktop — fires whenever files are added or changed -->
  <key>WatchPaths</key>
  <array>
    <string>/Users/PLACEHOLDER_USER/Desktop</string>
  </array>

  <!-- Also run once on load to catch anything already sitting on the Desktop -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/captures-watcher.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/captures-watcher-error.log</string>

  <key>LowPriorityIO</key>
  <true/>

  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
